# services/air.target
[Unit]
Description=Controls services required for an OpenHD Air unit

[Install]
WantedBy=multi-user.target
