# services/ground.target
[Unit]
Description=Controls services required for an OpenHD Ground unit

[Install]
WantedBy=multi-user.target
